# roles/khardix.https/files/dehydrated.timer
[Unit]
Description=Periodic renewal of HTTPS certificates

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=24hours

[Install]
WantedBy=timers.target
